# debian/postinst
#!/bin/sh
set -e

systemctl daemon-reload
systemctl enable --now dnsmasq-goui.service
